# system_files/etc/systemd/system/rpm-ostreed-automatic.service.d/override.conf
[Service]
ExecStartPost=-/usr/bin/sh -c 'rpm-ostree status --pending-exit-77 >/dev/null ; [ $? = 77 ] && touch /run/reboot-required'
